# cron/crontab
#crontab
@reboot root    *modules/get_study_index.py